# scripts/com.brandmint.twitter-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.brandmint.twitter-sync</string>

    <key>Comment</key>
    <string>Brandmint weekly Twitter sync — discovers prompts, downloads assets</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Volumes/madara/2026/brandmint/scripts/twitter_sync_all.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Volumes/madara/2026/brandmint/references/twitter-sync/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/madara/2026/brandmint/references/twitter-sync/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/sheshnarayaniyer</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Volumes/madara/2026/brandmint</string>
  </dict>
</plist>
